# launchd/com.copilot.watch-sessions.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.copilot.watch-sessions</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON3__</string>
        <string>tools/watch-sessions.py</string>
        <string>--daemon</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__HOME__/.copilot</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/.copilot/session-state/.watch-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.copilot/session-state/.watch-launchd.log</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
